# backend/posbah-backend.service
[Unit]
Description=PosBah Go Backend Service
After=network.target postgresql.service
# Jika backend crash karena env tidak ada, jangan restart tanpa batas
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=muizz9900
WorkingDirectory=/home/muizz9900
Environment=PORT=3000
# Muat semua secret dari file env (tidak di-commit ke git)
EnvironmentFile=/home/muizz9900/.posbah-env
ExecStart=/home/muizz9900/posbah-backend
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
